# misc/webhook-exec.service
[Unit]
Description=Webhook Exec

[Service]
ExecStart=/usr/local/bin/webhook-exec

[Install]
WantedBy=default.target
