# misc/webhook-exec.service
[Unit]
Description=Webhook Exec
After=network-online.target

[Service]
ExecStart=/usr/local/bin/webhook-exec

[Install]
WantedBy=default.target
